# tools/hotplug/Linux/systemd/xenstored-daemon.service.in
[Unit]
Description=The xenstore daemon in dom0
Requires=proc-xen.mount var-lib-xenstored.mount
After=proc-xen.mount var-lib-xenstored.mount
Before=libvirtd.service libvirt-guests.service
RefuseManualStop=true
ConditionPathExists=/proc/xen/capabilities

[Service]
PIDFile=@XEN_RUN_DIR@/xenstored.pid
ExecStartPre=/bin/grep -q control_d /proc/xen/capabilities
ExecStart=@XEN_SCRIPT_DIR@/launch-xenstore 'systemd' 'supervised-by-systemd'

[Install]
WantedBy=multi-user.target
Also=proc-xen.mount
Also=var-lib-xenstored.mount
